# roles/common/templates/etc/update-motd.d/19-ipaddresses
#!/bin/sh

###############################################################
#  Script     : 19-ipaddresses
#  Author     : Andre
#  Date       : 16/07/2024
#  Description: displays ipv4 and ipv6 addresses
#  Args       : none
###############################################################

printf "\n"

printf "IPv4 addresses\n"
printf "===============================================================================\n"
ip -h -br -f inet a 
printf "===============================================================================\n"
printf "IPv6 addresses\n"
ip -h -br -f inet6 a
printf "===============================================================================\n"